# DockerHelpers/unit3d.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:unit3d-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work
startsecs = 0
autostart=true
autorestart=true
user=www-data
numprocs=10

[program:unit3d-chat-server]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/node /usr/bin/laravel-echo-server start --dir=/var/www
autostart=true
autorestart=true
user=www-data
numprocs=1